# roles/common/files/profile.d/aliases.sh
# !!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!
# ! FILE IS CONTROLLED BY ANSIBLE, DO NOT CHANGE, OR ELSE YOUR CHANGES WILL BE EVENTUALLY LOST !
# !!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!
#
# Convenience aliases
#
alias dir="ls -F"
alias l="ls -l"
alias ll="ls -lF"
alias la="ls -la"
alias md=mkdir
alias ..="cd .."
alias ...="cd ../.."
alias motd="cat /etc/motd"
